# systemd/vectorcore-hss.service
[Unit]
Description=VectorCore HSS
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vectorcore
ExecStart=/opt/vectorcore/bin/hss -c /opt/vectorcore/etc/hss.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
